# rsyslog/shared-volume/app/supervisord.conf
[supervisord]
nodaemon = true
user = root
childlogdir = /var/log/custom

[program:sshd]
command=/usr/sbin/sshd -D -e
redirect_stderr = True
stdout_logfile = /var/log/custom/sshd.log
stderr_syslog = True
stdout_syslog = True
